# mysite.service
# systemd 서비스 파일
# 파일 위치: /etc/systemd/system/mysite.service
# 
# 설치 방법:
# 1. sudo cp mysite.service /etc/systemd/system/
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable mysite
# 4. sudo systemctl start mysite

[Unit]
Description=Django mysite Gunicorn Application Server
Documentation=https://docs.gunicorn.org/
After=network.target postgresql.service
Wants=postgresql.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
# 서비스 타입
Type=exec
Restart=on-failure
RestartSec=5

# 사용자 설정 (보안)
User=www-data
Group=www-data

# 환경 설정
Environment=DJANGO_SETTINGS_MODULE=config.settings.prod
Environment=PYTHONPATH=/home/ubuntu/projects/mysite
Environment=PYTHONUNBUFFERED=1

# 작업 디렉토리
WorkingDirectory=/home/ubuntu/projects/mysite

# 실행 명령
ExecStart=/home/ubuntu/projects/mysite/venv/bin/gunicorn \
    --config /home/ubuntu/projects/mysite/gunicorn.conf.py \
    config.wsgi:application
    
ExecReload=/bin/kill -s HUP $MAINPID

# 리소스 제한
LimitNOFILE=65536
LimitNPROC=4096

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/ubuntu/projects/mysite/logs /home/ubuntu/projects/mysite/media /home/ubuntu/projects/mysite/staticfiles

# 킬 시그널
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target